# .claude/worktrees/security-bugfixes-2026-06-02/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:backend]
command=/app/backend/server
directory=/app/backend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PORT="8080",AUTH_SECRET="%(ENV_AUTH_SECRET)s",DATABASE_URL="%(ENV_DATABASE_URL)s",DB_TYPE="%(ENV_DB_TYPE)s",ENV="%(ENV_ENV)s",REDIS_URL="%(ENV_REDIS_URL)s",ALLOWED_ORIGINS="%(ENV_ALLOWED_ORIGINS)s",NEXTAUTH_URL="%(ENV_NEXTAUTH_URL)s",NEXTAUTH_SECRET="%(ENV_NEXTAUTH_SECRET)s",OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s",ANTHROPIC_API_KEY="%(ENV_ANTHROPIC_API_KEY)s",NVIDIA_API_KEY="%(ENV_NVIDIA_API_KEY)s",GROQ_API_KEY="%(ENV_GROQ_API_KEY)s",GEMINI_API_KEY="%(ENV_GEMINI_API_KEY)s"

[program:frontend]
command=node /app/frontend/apps/web/server.js
directory=/app/frontend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=AUTH_SECRET="%(ENV_AUTH_SECRET)s",DATABASE_URL="%(ENV_DATABASE_URL)s",NEXTAUTH_URL="%(ENV_NEXTAUTH_URL)s",NEXTAUTH_SECRET="%(ENV_NEXTAUTH_SECRET)s",BACKEND_URL="%(ENV_BACKEND_URL)s",PORT="%(ENV_PORT)s",NODE_ENV="%(ENV_NODE_ENV)s"
